# deploy/gita-learning.service
[Unit]
Description=Gita Learning FastAPI app
After=network.target postgresql.service

[Service]
User=dev
Group=www-data
WorkingDirectory=/home/dev/gita-learning/backend
EnvironmentFile=/etc/gita-learning.env
ExecStart=/home/dev/gita-learning/backend/.venv/bin/gunicorn app.main:app -k uvicorn.workers.UvicornWorker --bind 127.0.0.1:8001 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
